# system/LaunchAgents/com.dotfiles.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles.update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__DOTFILES_DIR__/update.sh</string>
    </array>

    <!-- Run daily at 9 AM — edit Hour/Minute to change the schedule -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__DOTFILES_DIR__/logs/update.log</string>
    <key>StandardErrorPath</key>
    <string>__DOTFILES_DIR__/logs/update.log</string>

    <!-- Do not run immediately when loaded, only on the next scheduled time -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
